# profile.d/assh_alias.sh
#!/usr/bin/env bash

alias assh="ssh -o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null"
