# profile.d/hub.sh
#!/usr/bin/env bash

command -v hub &> /dev/null &&
  alias git="hub"
